# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10

[program:node-backend]
directory=/app/backend
command=pnpm start
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:certbot-renewal]
# Runs certbot renew every 12 hours and reloads nginx if renewed
command=/bin/bash -c "while true; do sleep 43200; certbot renew --nginx --quiet; done"
autostart=true
autorestart=true
priority=30